# files/etc/uci-defaults/90-thdn-reboot

#!/bin/sh
# 定时重启：每天 03:00 重启
echo "0 3 * * * reboot" >> /etc/crontabs/root
/etc/init.d/cron enable
/etc/init.d/cron start
exit 0
